# code/puppet/modules/devpro/templates/supervisord.conf.erb
[program:devpro]
command=/usr/local/bin/uwsgi
        --socket <%= run_dir %>/uwsgi.sock
        --chmod-socket
        --processes 2
        --master
        --virtualenv <%= venv_dir %>
        --pp <%= app_dir %>
        --module django.core.handlers.wsgi:WSGIHandler()
environment=DJANGO_SETTINGS_MODULE='devpro.settings'
user=root
autostart=true
autorestart=true
stdout_logfile=<%= log_dir %>/uwsgi.log
stderr_logfile=<%= log_dir %>/uwsgi.err
stopsignal=QUIT
